# bootup-slack-alert-service/bootup-slack-alert.service
[Unit]
Description=Slack alert on system boot

[Service]
Type=oneshot
ExecStart=/usr/local/bin/bootup-slack-alert.sh
