# bootup-slack-alert-service/bootup-slack-alert.timer
[Unit]
Description=Trigger slack alert 1 min after boot

[Timer]
OnBootSec=1min

[Install]
WantedBy=timers.target
